# backend/ai-dialogue-backend.service
[Unit]
Description=AI对话应用后端服务
Documentation=https://github.com/your-repo/ai-dialogue-backend
After=network.target
Wants=network.target

[Service]
Type=simple
User=backend
Group=backend
WorkingDirectory=/opt/ai-dialogue-backend
Environment=PATH=/opt/ai-dialogue-backend/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/ai-dialogue-backend
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# 基础配置
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=DEBUG=false
Environment=LOG_LEVEL=INFO

# LLM服务配置（请在部署时设置实际的API密钥）
Environment=OPENROUTER_API_KEY=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
Environment=OPENROUTER_BASE_URL=https://openrouter.ai/api/v1
Environment=OPENROUTER_MODEL=qwen/qwen3-235b-a22b:free

# STT服务配置
Environment=USE_REAL_VOSK=true
Environment=VOSK_MODEL_PATH=/opt/ai-dialogue-backend/model/vosk-model
Environment=VOSK_SAMPLE_RATE=16000

# 安全配置
Environment=ALLOWED_ORIGINS=["*"]

# 性能配置
Environment=STT_TIMEOUT=20
Environment=LLM_TIMEOUT=45
Environment=WEBSOCKET_TIMEOUT=600
Environment=MAX_WORKERS=8

# 日志配置
Environment=LOG_FILE=/var/log/ai-dialogue-backend/app.log
Environment=LOG_MAX_SIZE=100MB
Environment=LOG_BACKUP_COUNT=10
Environment=LOG_FORMAT=json

# 启动命令
ExecStart=/opt/ai-dialogue-backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
TimeoutStopSec=30

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ai-dialogue-backend /opt/ai-dialogue-backend/logs
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
LockPersonality=true

# 网络安全
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target